# autoinstall-usb/overlay/etc/init.d/autoprov-mount
#!/sbin/openrc-run

description="Mount ENV/PERSIST and prepare /run/autoprov.env"

depend() {
  need localmount
  before autoprov
}

start() {
  ebegin "Mounting ENV/PERSIST"
  /opt/autorun/mount-env.sh
  eend $?
}
